# openpryv/openpryv.service
[Unit]
Description=Open-Pryv.io Service

[Service]
Type=simple
ExecStart=/bin/bash /usr/bin/openpryv.sh

[Install]
WantedBy=multi-user.target